# clawboard.service
[Unit]
Description=Clawboard Dashboard
After=network.target

[Service]
Type=simple
User=claw
WorkingDirectory=/home/claw/.openclaw/workspace/clawboard
ExecStart=/home/claw/.openclaw/workspace/clawboard/clawboard
EnvironmentFile=/home/claw/.openclaw/workspace/clawboard/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
